# conf/supervisord-lamp.conf
[supervisord]
stdout_logfile=/var/log/supervisor/supervisor.log
nodaemon=true

[program:mysqld]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true

[program:apache2]
command=/start-apache2.sh
numprocs=1
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
